# config/systemd/.config/systemd/user/check-temp-trends.service
[Unit]
Description=Check long-term CPU/GPU temperature trends

[Service]
Type=oneshot
ExecStart=%h/.config/check-temp-trends.sh

