# config/systemd/.config/systemd/user/dns-check.service
[Unit]
Description=Check DNS configuration

[Service]
Type=oneshot
ExecStart=%h/.config/check-dns.sh
