# config/systemd/.config/systemd/user/log-temps.service
[Unit]
Description=Log CPU/GPU temperatures

[Service]
ExecStart=%h/.config/log-temps.sh
