# config/systemd/.config/systemd/user/update-check.service
[Unit]
Description=Check for system updates and age

[Service]
ExecStart=%h/.config/update-timer-task.sh

[Install]
WantedBy=default.target
