# config/systemd/.config/systemd/user/update-check.timer
[Unit]
Description=Run update check every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
